# lib/systemd/bridgehead-update@.timer
[Unit]
Description=Nightly Updates of Bridgehead (%i)

[Timer]
OnCalendar=*-*-* 03:00:00

[Install]
WantedBy=basic.target
