# lib/systemd/bridgehead-update@.service
<<<<<<< HEAD
[Unit]
Description=Bridgehead (%i) Update Service

[Service]
Type=oneshot
WorkingDirectory=/srv/docker/bridgehead/
ExecStart=/bin/bash -c "/srv/docker/bridgehead/update-bridgehead.sh %i"

[Install]
WantedBy=multi-user.target
=======
[Unit]
Description=Bridgehead (%i) Update Service

[Service]
Type=oneshot
User=bridgehead
ExecStart=/srv/docker/bridgehead/bridgehead update %i

[Install]
WantedBy=multi-user.target
>>>>>>> version-1
