# lib/systemd/bridgehead-update@.service
[Unit]
Description=Bridgehead (%i) Update Service

[Service]
Type=oneshot
User=bridgehead
ExecStart=/srv/docker/bridgehead/bridgehead update %i

[Install]
WantedBy=multi-user.target
