# lib/systemd/bridgehead@.service
[Unit]
Description=Bridgehead (%i) Service

[Service]
User=bridgehead
Restart=always
RestartSec=30
ExecStart=/srv/docker/bridgehead/bridgehead start %i
ExecStop=/srv/docker/bridgehead/bridgehead stop %i

[Install]
WantedBy=multi-user.target
