# lib/systemd/bridgehead-update@.timer
[Unit]
Description=Hourly Updates of Bridgehead (%i)

[Timer]
OnCalendar=*-*-* *:00:00

[Install]
WantedBy=basic.target
