# lib/systemd/bridgehead@.service
[Unit]
Description=Bridgehead (%i) Service

[Service]
User=bridgehead
Restart=always
RestartSec=30

WorkingDirectory=/srv/docker/bridgehead/

ExecStart=/srv/docker/bridgehead/start-bridgehead.sh %i
RemainAfterExit=true
ExecStop=/srv/docker/bridgehead/stop-bridgehead.sh %i

[Install]
WantedBy=multi-user.target
