# convenience/bridgehead@.service
[Unit]
Description=Bridgehead (%i) Service

[Service]
Restart=always
RestartSec=30

WorkingDirectory=/srv/docker/bridgehead/%i

ExecStartPre=docker-compose --env-file ../site-config/%i.env down
ExecStart=docker-compose --env-file ../site-config/%i.env up
ExecStop=docker-compose --env-file ../site-config/%i.env down

[Install]
WantedBy=multi-user.target
